# launchd/com.user.nextdns-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.nextdns-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/usr/local/bin/sync_nextdns_local.py</string>
    </array>

    <!-- Run every 30 minutes -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- Load env var from ~/.zshenv at launch time -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>NEXTDNS_API_KEY</key>
        <string>REPLACE_WITH_YOUR_API_KEY</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/nextdns-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/nextdns-sync.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
